# service/tuerd/run
#!/bin/sh

exec 2>&1

export TUERD_DEBUG=0
export TUERD_REPO_PATH="/home/tuerd/git"
export USBLOCKD_PIDFILE="/var/run/usblockd.pid"

echo $$ > /var/run/tuerd.pid

sleep 5

exec chpst -utuerd:tuerd:dialout /usr/local/sbin/tuerd
